# deploy/systemd/radio-proxies.service
[Unit]
Description=Radio Dashboard — 8 ShazamIO proxies
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%h/dev/radio-playlist-dashboard
ExecStart=%h/dev/radio-playlist-dashboard/.venv/bin/python scripts/proxy_manager.py start
ExecStop=%h/dev/radio-playlist-dashboard/.venv/bin/python scripts/proxy_manager.py stop
TimeoutStartSec=120

[Install]
WantedBy=default.target
