# scripts/1036updater.service
[Unit]
Description=Radio Playlist Dashboard — Collector Daemon
Documentation=https://github.com/brchn6/radio-playlist-dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/dev/radio-playlist-dashboard
ExecStart=/usr/bin/env python %h/dev/radio-playlist-dashboard/scripts/updater.py
Restart=always
RestartSec=10

# NO GIT_AUTO_PUSH. The collector must never push to git — that was the
# ~720-commits/day pattern that risked a ToS strike. Data goes to Supabase.
# Credentials are read from the project's .env (SUPABASE_URL / SUPABASE_SECRET_KEY).
Environment=RETENTION_DAYS=45

[Install]
WantedBy=default.target

# ── Not installed by default ─────────────────────────────────────────────
# The collector currently runs under `nohup` via scripts/manage.sh, so it does
# NOT survive a reboot. To fix that:
#
#   systemctl --user enable --now 1036updater.service
#   loginctl enable-linger barc      # keeps user units running after logout
#
# Do this only if you also move the proxies to a unit — updater.py is useless
# without them. See .planning/DEPLOY-ARCHITECTURE.md "Known limits".
